# ops/systemd/star-sign-backup.timer
[Unit]
Description=Star Sign PostgreSQL backup — daily 02:30 UTC

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
